# package/utils/busybox/files/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:/bin/cttyhack /bin/sh
